# valmindiaqueue.service
[Unit]
Description=VEGA Assassins Matchmaking Discord Bot
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/VEGA-Queue-System
ExecStart=/home/pi/VEGA-Queue-System/venv/bin/python /home/pi/VEGA-Queue-System/bot.py
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/home/pi/VEGA-Queue-System/logs/bot.log
StandardError=append:/home/pi/VEGA-Queue-System/logs/bot.error.log

[Install]
WantedBy=multi-user.target


